# root/etc/services.d/ombi/run
#!/usr/bin/with-contenv bash

cd /app/Ombi
exec s6-setuidgid ombi mono Ombi.exe "${RUN_OPTS}"
